# docker/blah-mesos-slave/supervisor/filebeat.conf
[program:filebeat]
command=/opt/filebeat/filebeat -e -c /filebeat.yml
